# .config/autostart/xmodmap.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=xmodmap
Exec=/home/tanishiking/dotfiles/bin/exec_xmodmap
Comment=
